# com.ben.session-mirror.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.ben.session-mirror</string>
  <key>ProgramArguments</key><array>
    <string>/usr/bin/python3</string>
    <string>/Users/ben/daily/bin/corpus</string>
    <string>mirror-sessions</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ProcessType</key><string>Background</string>
  <key>StandardOutPath</key><string>/Users/ben/Library/Logs/daily-note/session-mirror.log</string>
  <key>StandardErrorPath</key><string>/Users/ben/Library/Logs/daily-note/session-mirror.log</string>
</dict></plist>
