# com.ben.granola-reconcile.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.ben.granola-reconcile</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/ben/code/daily-auto/granola-reconcile.sh</string>
	</array>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>ExitTimeOut</key>
	<integer>300</integer>

	<key>Nice</key>
	<integer>10</integer>

	<key>ProcessType</key>
	<string>Background</string>

	<key>StandardOutPath</key>
	<string>/Users/ben/Library/Logs/daily-note/launchd-granola-stdout.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/ben/Library/Logs/daily-note/launchd-granola-stderr.log</string>
</dict>
</plist>
